# system/units/appscale-control.target
[Unit]
Description=AppScale Management/Controllers
Requires=appscale-routing.target
After=appscale-routing.target
Conflicts=appscale-down.target