# system/units/appscale.target
[Unit]
Description=AppScale GTS
Wants=appscale-control.target
After=appscale-control.target
Conflicts=appscale-down.target
