# system/units/appscale-service.target
[Unit]
Description=AppScale Intermediate Services
Requires=appscale-backend.target
After=appscale-backend.target
Conflicts=appscale-down.target