# system/units/appscale-down.target
[Unit]
Description=Shutdown AppScale GTS services