# system/units/appscale-routing.target
[Unit]
Description=AppScale routing
Requires=appscale-instance.target
After=appscale-instance.target
Conflicts=appscale-down.target