# system/units/appscale-backend.target
[Unit]
Description=AppScale GTS backend/foundational services
Before=appscale-down.target
Conflicts=appscale-down.target