# system/units/appscale-instance.target
[Unit]
Description=AppScale Application Instances
Requires=appscale-service.target
After=appscale-service.target
Conflicts=appscale-down.target