# system/units.d/ejabberd.service.d/10-appscale-ejabberd.conf
# Unit configuration for ejabberd, /lib/systemd/system/ejabberd.service.d/appscale-ejabberd.conf

[Unit]
PartOf=appscale-ejabberd.target
ReloadPropagatedFrom=appscale-ejabberd.target